# scripts/supervisor-disque.conf
[supervisord]
nodaemon=true

[program:disque]
command=disque-server
